# fs/etc/services.d/alertmanager/run
#!/usr/bin/with-contenv ash

cd /alertmanager

exec runas /bin/alertmanager --config.file=/config/alertmanager.yml --storage.path=/alertmanager $ALERTMANAGER_ARGS